# templates/default/redis.upstart.conf.erb
description "Redis Server"
author "Edward Paget"

start on [2345]
stop on [!12345]

expect fork

respawn
respawn limit 10 5

setuid <%= @user %>
setgid <%= @group %>

exec /usr/local/bin/redis-server /etc/redis/<%= @port %>.conf